# batch_task/cron.d/crontab
*/5 * * * * python3 -c 'print("CRON RUN")'
